# files/default/xvfb.conf
description "xvfb"
start on runlevel [2345]
exec /usr/bin/Xvfb :99
